# systemd/flowless-watchdog.service
[Unit]
Description=Flowless Watchdog - Auto-restart crashed backends
Documentation=https://github.com/ArashDoDo2/flowless
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Note: Running as root is required to use systemctl restart
# For production, consider using sudo with limited permissions
User=root
ExecStart=/usr/local/bin/flowless-watchdog
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=flowless-watchdog

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict

# Allow watchdog to restart services
# These capabilities are needed for systemctl operations
CapabilityBoundingSet=CAP_KILL CAP_SYS_ADMIN

[Install]
WantedBy=multi-user.target
